# launchd/com.laodai.codex-weekly-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.laodai.codex-weekly-backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/daizhengzhou/Documents/老戴Ai备份外脑/scripts/backup_and_sync_codex.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/daizhengzhou/Documents/老戴Ai备份外脑/backups/codex-weekly/weekly-backup-and-github-sync.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/daizhengzhou/Documents/老戴Ai备份外脑/backups/codex-weekly/weekly-backup-and-github-sync.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
